# src/modules/fullpageos/filesystem/root_init/etc/systemd/system/enable_gpu_first_boot.service
[Unit]
Description=on first boot set up GPU acceleration
[Service]
ExecStart=/home/pi/scripts/enable_gpu
ExecStop=
[Install]
WantedBy=multi-user.target
